# timescaledb/rootfs/etc/s6-overlay/s6-rc.d/postgres/run
#!/command/with-contenv bashio
exec gosu postgres postgres -D /data/postgres -k /tmp
